# system_cfg/etc/systemd/system/postgresql.service.d/override.conf
[Unit]
After=docker.service
After=network.target

[Service]
ExecStartPre=/usr/bin/sleep 20